# ansible/roles/llm-mlx/templates/com.hydra.mlx-server.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  /Library/LaunchDaemons/com.hydra.mlx-server.plist
  Managed by Ansible (llm-mlx role) — do not edit manually.
  Apple MLX inference server (mlx_lm.server).
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hydra.mlx-server</string>

  <key>ProgramArguments</key>
  <array>
{% if mlx_server_backend == 'vllm_mlx' %}
    <string>{{ mlx_venv }}/bin/vllm-mlx</string>
    <string>serve</string>
    <string>{{ mlx_default_model }}</string>
    <string>--host</string><string>{{ mlx_engine_bind_host }}</string>
    <string>--port</string><string>{{ mlx_engine_bind_port }}</string>
    <string>--enable-metrics</string>
    <string>--continuous-batching</string>
    <string>--use-paged-cache</string>
    <string>--max-request-tokens</string><string>{{ mlx_max_tokens }}</string>
    <string>--cache-memory-percent</string><string>0.20</string>
{% else %}
    {# Pure mlx_lm.server has no Prometheus metrics; the perf-proxy sidecar
       (com.hydra.mlx-perf) fronts it on {{ mlx_port }} and the engine binds the
       loopback port {{ mlx_engine_bind_port }}. #}
    <string>{{ mlx_venv }}/bin/python</string>
    <string>-m</string><string>mlx_lm.server</string>
    <string>--model</string><string>{{ mlx_default_model }}</string>
    <string>--host</string><string>{{ mlx_engine_bind_host }}</string>
    <string>--port</string><string>{{ mlx_engine_bind_port }}</string>
    <string>--max-tokens</string><string>{{ mlx_max_tokens }}</string>
{% endif %}
  </array>

  <key>WorkingDirectory</key>
  <string>{{ mlx_dir }}</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>TOKENIZERS_PARALLELISM</key>
    <string>false</string>
    <key>HF_DATASETS_OFFLINE</key>
    <string>1</string>
    <key>TRANSFORMERS_OFFLINE</key>
    <string>1</string>
    <key>HF_HUB_OFFLINE</key>
    <string>1</string>
    <key>MLX_METAL_DEBUG</key>
    <string>0</string>
    <key>HYDRA_MLX_MAX_KV_SIZE</key>
    <string>{{ mlx_max_kv_size }}</string>
    <!-- LLM-6: OTLP trace export to the local OTEL agent → gateway -->
    <key>OTEL_TRACES_EXPORTER</key>
    <string>otlp</string>
    <key>OTEL_EXPORTER_OTLP_TRACES_ENDPOINT</key>
    <string>{{ otel_traces_http_endpoint | default('http://127.0.0.1:4318') }}/v1/traces</string>
    <key>OTEL_SERVICE_NAME</key>
    <string>mlx-server</string>
    <key>OTEL_RESOURCE_ATTRIBUTES</key>
    <string>node_id={{ hydra_node_id | default(inventory_hostname) }},component=mlx</string>
  </dict>

  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>ThrottleInterval</key><integer>15</integer>

  <key>StandardOutPath</key>
  <string>/var/log/hydra/mlx-server.log</string>
  <key>StandardErrorPath</key>
  <string>/var/log/hydra/mlx-server-error.log</string>
</dict>
</plist>
