# monitoring/otel-agent-linux/otelcol-agent.service
[Unit]
Description=OpenTelemetry Collector Agent (Hydra dev node)
Documentation=https://opentelemetry.io
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/otelcol-contrib --config /etc/otel/agent-config.yaml
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=otelcol-agent

# Safety limits
LimitNOFILE=65536
MemoryMax=256M
CPUQuota=20%

[Install]
WantedBy=multi-user.target
