# apple-silicon-monitoring/deploy/systemd/otel-gateway.service
[Unit]
Description=OpenTelemetry Collector Gateway
Documentation=https://opentelemetry.io/docs/collector/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=otel
Group=otel

# Configuration
EnvironmentFile=/etc/otel-collector/gateway.env
ExecStart=/usr/local/bin/otelcol-contrib --config=/etc/otel-collector/gateway-config.yaml
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Resource limits
MemoryMax=3G
MemoryHigh=2.5G
CPUQuota=400%

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/otel-gateway /var/log/otel-gateway

# Capabilities
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=otel-gateway

[Install]
WantedBy=multi-user.target
